# scripts/launchagents/com.dailystock.moomoo-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dailystock.moomoo-sync</string>

    <!-- Pull Moomoo trade history every 15 minutes. SIMULATE by default;
         flip to LIVE by editing MOOMOO_TRADE_ENV in the project .env -->
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__PROJECT_DIR__/scripts/sync_moomoo_live.py</string>
        <string>--days</string>
        <string>2</string>
        <string>-v</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Run on login + every 900 seconds (15 min) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchagent.moomoo-sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchagent.moomoo-sync.err.log</string>
</dict>
</plist>
